# config/crontab
* * * * * php /app/artisan schedule:run > /dev/null &