# config/crontab
* * * * * /usr/bin/php /www/artisan schedule:run > /dev/null &